# qr3_94.conf
[program:qr3_94]
command=python3 /home/curator3/qr3_94/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_94/test.err.log
stdout_logfile=/home/curator3/qr3_94/test.err.log
